# profile.d/10-hooks.zsh
# shellcheck shell=bash

zsh_bashcompinit