# external/profile.d/jetbrains.sh
# shellcheck shell=sh
export AQUA_PROPERTIES="/Users/j5pu/JetBrains/Aqua/idea.properties"
export AQUA_VM_OPTIONS="/Users/j5pu/JetBrains/Aqua/aqua.vmoptions"
export CLION_PROPERTIES="/Users/j5pu/JetBrains/CLion/idea.properties"
export CLION_VM_OPTIONS="/Users/j5pu/JetBrains/CLion/clion.vmoptions"
export CLIENT_PROPERTIES="/Users/j5pu/JetBrains/Client/idea.properties"
export CLIENT_VM_OPTIONS="/Users/j5pu/JetBrains/Client/client.vmoptions"
export DATAGRIP_PROPERTIES="/Users/j5pu/JetBrains/DataGrip/idea.properties"
export DATAGRIP_VM_OPTIONS="/Users/j5pu/JetBrains/DataGrip/datagrip.vmoptions"
export DATASPELL_PROPERTIES="/Users/j5pu/JetBrains/DataSpell/idea.properties"
export DATASPELL_VM_OPTIONS="/Users/j5pu/JetBrains/DataSpell/dataspell.vmoptions"
export GATEWAY_PROPERTIES="/Users/j5pu/JetBrains/Gateway/idea.properties"
export GATEWAY_VM_OPTIONS="/Users/j5pu/JetBrains/Gateway/gateway.vmoptions"
export GOLAND_PROPERTIES="/Users/j5pu/JetBrains/GoLand/idea.properties"
export GOLAND_VM_OPTIONS="/Users/j5pu/JetBrains/GoLand/goland.vmoptions"
export IDEA_PROPERTIES="/Users/j5pu/JetBrains/Idea/idea.properties"
export IDEA_VM_OPTIONS="/Users/j5pu/JetBrains/Idea/idea.vmoptions"
export PYCHARM_PROPERTIES="/Users/j5pu/JetBrains/PyCharm/idea.properties"
export PYCHARM_VM_OPTIONS="/Users/j5pu/JetBrains/PyCharm/pycharm.vmoptions"
export RUBYMINE_PROPERTIES="/Users/j5pu/JetBrains/RubyMine/idea.properties"
export RUBYMINE_VM_OPTIONS="/Users/j5pu/JetBrains/RubyMine/rubymine.vmoptions"
export STUDIO_PROPERTIES="/Users/j5pu/JetBrains/Studio/idea.properties"
export STUDIO_VM_OPTIONS="/Users/j5pu/JetBrains/Studio/studio.vmoptions"
export WEBSTORM_PROPERTIES="/Users/j5pu/JetBrains/WebStorm/idea.properties"
export WEBSTORM_VM_OPTIONS="/Users/j5pu/JetBrains/WebStorm/webstorm.vmoptions"
