# external/profile.d/jetbrains.sh
# shellcheck shell=sh
export APPCODE_PROPERTIES="/Users/j5pu/JetBrains/bin/AppCode/idea.properties"
export APPCODE_VM_OPTIONS="/Users/j5pu/JetBrains/bin/AppCode/appcode.vmoptions"
export AQUA_PROPERTIES="/Users/j5pu/JetBrains/bin/Aqua/idea.properties"
export AQUA_VM_OPTIONS="/Users/j5pu/JetBrains/bin/Aqua/aqua.vmoptions"
export CLION_PROPERTIES="/Users/j5pu/JetBrains/bin/CLion/idea.properties"
export CLION_VM_OPTIONS="/Users/j5pu/JetBrains/bin/CLion/clion.vmoptions"
export DATAGRIP_PROPERTIES="/Users/j5pu/JetBrains/bin/DataGrip/idea.properties"
export DATAGRIP_VM_OPTIONS="/Users/j5pu/JetBrains/bin/DataGrip/datagrip.vmoptions"
export DATASPELL_PROPERTIES="/Users/j5pu/JetBrains/bin/DataSpell/idea.properties"
export DATASPELL_VM_OPTIONS="/Users/j5pu/JetBrains/bin/DataSpell/dataspell.vmoptions"
export GATEWAY_PROPERTIES="/Users/j5pu/JetBrains/bin/Gateway/idea.properties"
export GATEWAY_VM_OPTIONS="/Users/j5pu/JetBrains/bin/Gateway/gateway.vmoptions"
export GOLAND_PROPERTIES="/Users/j5pu/JetBrains/bin/GoLand/idea.properties"
export GOLAND_VM_OPTIONS="/Users/j5pu/JetBrains/bin/GoLand/goland.vmoptions"
export IDEA_PROPERTIES="/Users/j5pu/JetBrains/bin/Idea/idea.properties"
export IDEA_VM_OPTIONS="/Users/j5pu/JetBrains/bin/Idea/idea.vmoptions"
export PYCHARM_PROPERTIES="/Users/j5pu/JetBrains/bin/PyCharm/idea.properties"
export PYCHARM_VM_OPTIONS="/Users/j5pu/JetBrains/bin/PyCharm/pycharm.vmoptions"
export RUBYMINE_PROPERTIES="/Users/j5pu/JetBrains/bin/RubyMine/idea.properties"
export RUBYMINE_VM_OPTIONS="/Users/j5pu/JetBrains/bin/RubyMine/rubymine.vmoptions"
export WEBSTORM_PROPERTIES="/Users/j5pu/JetBrains/bin/WebStorm/idea.properties"
export WEBSTORM_VM_OPTIONS="/Users/j5pu/JetBrains/bin/WebStorm/webstorm.vmoptions"
